# deployment/ansible-playbook/bisq-daonode.service.j2
[Unit]
Description=daonode

[Service]
ExecStart=/home/user/run_daonode.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target